# squashfs-root/usr/lib/opkg/info/rpcd.list
/etc/init.d/rpcd
/etc/config/rpcd
/sbin/rpcd
