# squashfs-root/usr/lib/opkg/info/sambuca-web-device.list
/etc/uci-defaults/52_ftp
/etc/conndev
/etc/uci-defaults/57_hosts
/usr/libexec/rpcd/vlan
/www/langs/es-es.json
/etc/channels/5Glist
/usr/release_notes.txt
/www/langs/zh-cn.json
/usr/libexec/rpcd/wps
/www/favicons/favicon-16x16.png
/www/langs/id.json
/www/favicons/apple-touch-icon-144x144.png
/etc/uci-defaults/62_dhcp
/www/favicons/apple-touch-icon-76x76.png
/etc/uci-defaults/54_sysupgrade
/www/build/geogtq-bd-webfont.woff
/etc/wps_connection
/lib/network/hosts.sh
/usr/libexec/rpcd/wifi
/www/build/geog-thit-webfont.woff
/www/langs/fr-fr.json
/www/build/android-illustration.png
/etc/dmz_ordering
/www/langs/es-ar.json
/www/build/app.js
/usr/libexec/rpcd/chlist
/etc/uci-defaults/55_redirect
/usr/libexec/rpcd/connection
/www/build/windows-illustration.png
/www/favicons/android-chrome-144x144.png
/www/build/geog-bdit-webfont.woff
/usr/libexec/rpcd/ifconn
/www/build/main.css
/etc/init.d/fwupdate
/www/langs/nl.json
/usr/libexec/rpcd/ntpd
/etc/uci-defaults/99_xbox_sku
/www/langs/de.json
/www/langs/fi.json
/usr/libexec/rpcd/pppoe
/etc/uci-defaults/60_usbconf
/www/langs/vi.json
/www/favicons/android-chrome-36x36.png
/www/favicons/favicon-96x96.png
/www/favicons/apple-touch-icon-180x180.png
/etc/uci-defaults/61_disable_wan_pings
/www/langs/tr.json
/www/langs/it-it.json
/www/build/modemTroubleshooting.svg
/usr/libexec/rpcd/attdev
/etc/hotplug.d/iface/12-ip_check
/etc/uci-defaults/59_sambuca
/sbin/usb-detect
/www/build/mac-illustration.png
/etc/init.d/connled
/usr/libexec/rpcd/guestwifi
/etc/init.d/wan-detect
/www/favicons/favicon-64x64.png
/etc/firewall.eula
/www/favicons/android-chrome-48x48.png
/etc/uci-defaults/56_led
/www/langs/ru.json
/usr/libexec/rpcd/sambuca
/usr/share/rpcd/acl.d/superuser.json
/www/build/wrt@2x.png
/www/build/killer-text-only.png
/usr/libexec/rpcd/email
/www/langs/ar-sa.json
/www/build/geog-ulit-webfont.woff
/www/favicons/favicon-32x32.png
/etc/uci-defaults/58_wireless
/sbin/connled
/www/cgi-bin/binaryupload
/www/langs/el.json
/www/favicons/android-chrome-72x72.png
/etc/channels/2Glist
/etc/uci-defaults/40_eula
/usr/libexec/rpcd/fwupdate
/www/favicons/apple-touch-icon-60x60.png
/www/build/unplugModem.svg
/www/build/iphone-illustration.png
/www/favicons/apple-touch-icon-57x57.png
/www/build/geogtq-th-webfont.woff
/www/build/geogrotesque-regular-webfont.woff
/www/build/windows-illustration.svg
/www/langs/pt-pt.json
/www/favicons/android-chrome-192x192.png
/www/build/geog-rgit-webfont.woff
/www/build/ipad-illustration.png
/www/langs/pl.json
/www/build/geogrotesque-medium-webfont.woff
/usr/libexec/rpcd/ookla
/www/langs/da.json
/etc/sysupgrade.conf
/www/langs/pt-br.json
/usr/libexec/rpcd/wizupdate
/www/build/geogtq-ul-webfont.woff
/www/langs/nb.json
/usr/libexec/rpcd/ipv6
/www/favicons/apple-touch-icon-120x120.png
/www/favicons/apple-touch-icon-72x72.png
/www/cgi-bin/venomping
/etc/init.d/usb-detect
/www/build/geog-mdit-webfont.woff
/www/build/geogtq-sb-webfont.woff
/www/favicons/apple-touch-icon-114x114.png
/www/favicons/apple-touch-icon-152x152.png
/www/langs/th.json
/www/build/confirm@2x.png
/www/favicons/android-chrome-96x96.png
/etc/hotplug.d/usb/10-usb
/www/favicons/favicon-24x24.png
/www/build/killer-texture.png
/www/langs/sv.json
/etc/wordlist.txt
/usr/libexec/rpcd/cpu
/www/build/geog-sbit-webfont.woff
/etc/uci-defaults/53_disable_services
/www/build/mac-illustration.svg
/etc/uci-defaults/50_rootpass
/etc/openvpn/openvpn_updown
/etc/hotplug.d/scsi/11-scsi
/www/langs/timezones-en.json
/www/langs/fr-ca.json
/www/langs/en.json
/usr/libexec/rpcd/admin
/usr/share/rpcd/acl.d/unauthenticated.json
/etc/init.d/pppoe-detect
/www/index.html
/usr/libexec/rpcd/led
/etc/uci-defaults/51_samba
/usr/libexec/rpcd/usb
