# squashfs-root/usr/lib/opkg/info/siproxd.list
/etc/config/siproxd
/etc/init.d/siproxd
/usr/sbin/siproxd
